# chess-ai.service
[Unit]
Description=AI Chess — Llama vs Llama Flask App
After=network.target ollama.service
Wants=ollama.service

[Service]
Type=simple
User=<Your Name>
WorkingDirectory=/path/to/Chess-AI
ExecStart=/path/to/Chess-AI-main/start_chess.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
